# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

uci set network.lan.stp=1
uci commit network
/etc/init.d/network reload
luci-reload

exit 0
